# crontab
00 02 * * * archvsync /script.sh
#Una linea de regalo para que funcione
